# tools/docker/www/crontab
*/30 * * * * goaccess /var/log/caddy/access.log -o /app/templates/report.html --log-format=CADDY --ignore-crawlers 2>&1
