# tools/docker/www/crontab
*/30 * * * * goaccess /var/log/nginx/access.log -o /var/www/templates/report.html --log-format=COMBINED --ignore-crawlers 2>&1
